# templates/etc/systemd/system/docker.service.d/proxy.conf.j2
[Service]    
Environment="HTTP_PROXY={{ proxy_env.http_proxy }}" "NO_PROXY=localhost,127.0.0.1,swsnet.ch,swisssign.net"
